# conpaas-mysql/scripts/etc/supervisor/conf.d/mysql.conf
[program:mysqld]
command=/usr/sbin/mysqld
autostart=true
autorestart=true